# supervisord.conf
[unix_http_server]          ; supervisord 的 unix socket 服务配置
file=/tmp/supervisor.sock   ; socket 文件的保存目录


[supervisord]                ; supervisord 的主进程配置
logfile=/tmp/supervisord.log ; 主要的进程日志配置
logfile_maxbytes=50MB        ; 最大日志体积, 默认 50MB
logfile_backups=10           ; 日志文件备份数目, 默认 10
loglevel=info                ; 日志级别, 默认 info; 还有: debug, warn, trace
pidfile=/tmp/supervisord.pid ; supervisord 的 pidfile 文件
nodaemon=false               ; 是否以守护进程的方式启动
minfds=1024                  ; 最小的有效文件描述符, 默认 1024
minprocs=200                 ; 最小的有效进程描述符, 默认 200


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket


[program:bet_server]
command=python3 -u /root/bet-server/src/app.py
autorestart=true
directory=/root/bet-server
stdout_logfile=/root/bet-server/bet_stdout.log
stderr_logfile=/root/bet-server/bet_stderr.log


[program:sync_event]
command=python3 -u /root/bet-server/src/sync_event.py
autorestart=true
directory=/root/bet-server
stdout_logfile=/root/bet-server/sync_stdout.log
stderr_logfile=/root/bet-server/sync_stderr.log
